# scripts/com.polymarket-agent.autotune.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polymarket-agent.autotune</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/pharrelly/codebase/github/polymarket-agent/scripts/autotune.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>WorkingDirectory</key>
    <string>/Users/pharrelly/codebase/github/polymarket-agent</string>

    <key>StandardOutPath</key>
    <string>/Users/pharrelly/codebase/github/polymarket-agent/logs/autotune/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/pharrelly/codebase/github/polymarket-agent/logs/autotune/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <!-- Provider: "claude" (default), "openai", or "anthropic" -->
        <!-- <key>AUTOTUNE_PROVIDER</key> -->
        <!-- <string>openai</string> -->
        <!-- <key>AUTOTUNE_MODEL</key> -->
        <!-- <string>gpt-4o</string> -->
        <!-- Set the API key for your chosen provider -->
        <!-- <key>ANTHROPIC_API_KEY</key> -->
        <!-- <string>sk-ant-...</string> -->
        <!-- <key>OPENAI_API_KEY</key> -->
        <!-- <string>sk-...</string> -->
    </dict>
</dict>
</plist>
